# rootfs/etc/services.d/cloudflare-warp/run
#!/command/with-contenv bash
echo "[$(basename "$(pwd)")] starting..."
set -e
s6-svwait -U /run/service/d-bus
rm -f /run/cloudflare-warp/warp_service
exec 2>&1 s6-notifyoncheck /usr/bin/warp-svc
