# cloudflare-warp-supervisord.conf
[program:cloudflare-warp-daemon]
autostart=0
command=/usr/bin/warp-svc

[program:cloudflare-warp-startup]
autorestart=0
startsecs=0
command=/usr/lib/supervisor/scripts/cloudflare-warp-startup.sh
